# templates/storyboard-workers.conf.erb
description "StoryBoard Deferred Workers"

start on started network-services
stop on stopping network-services

setuid <%= @user %>
setgid <%= @group %>

script
  <%= @install_root_api %>/bin/storyboard-worker-daemon --worker-count <%= @worker_count %>
end script